# etc/cron.d/myproject-snapshot-cassandra-storage
# Creates a cassandra (and btrfs) snapshot on the parser every hour.
#
0   *  *   *   *   root  /bin/bash -c '[[ "$(cat /.env|grep RUNS_PARSER=true)" != "" ]] && /root/scripts/snapshot-cassandra-storage.sh' >/dev/null 2>&1
