# etc/cron.d/myproject-back-up-mongodb
# Backs up the most recent btrfs snapshot of the mongodb data store to S3 every day at 12 AM CET (3 AM Oregon).
#
0   12 *   *   *   root  /bin/bash -c '[[ "$(cat /.env|grep RUNS_MONGOBACKUP=true)" != "" ]] && /root/scripts/back-up-mongodb.sh' >/dev/null 2>&1
