# etc/cron.d/myproject-pull-server-conf
# Looks for changes in the server configuration repo and pulls them, every 3-6 minutes
#
*/3 *  *   *   *   root  /bin/sleep $(shuf -i 0-90 -n 1)s; /bin/bash -c "(/usr/bin/flock -n 10 || exit 0; /etc/init.d/myproject-pull-server-conf start; rm /var/lock/myproject-pull-server-conf) 10>/var/lock/myproject-pull-server-conf" >/dev/null 2>&1
