# etc/cron.d/myproject-back-up-cassandra
# Backs up the most recent btrfs snapshot of the cassandra data store to S3 every day at 10 AM UTC (3 AM Oregon).
#
0   10 *   *   *   root  /bin/bash -c '[[ "$(cat /.env|grep RUNS_CASBACKUP=true)" != "" ]] && /root/scripts/back-up-cassandra.sh' >/dev/null 2>&1
